# systemd/system/slimlock.service
[Unit]
Description=Lock after sleep

[Service]
User=kciredor
ExecStart=/usr/bin/slimlock

[Install]
WantedBy=sleep.target
